# scripts/timer_test.service
[Unit]
Description=timer_test service

[Service]
ExecStart=/bin/bash /opt/scripts/timer_test.sh
